# _settings/supervisor/gunicorn.conf
[program:gunicorn]
command=/home/archi/PycharmProjects/sms/env/bin/gunicorn sms.wsgi
directory=/home/archi/PycharmProjects/sms/src
autostart=false
autorestart=true
