# _settings/supervisor/uwsgi.conf
[program:uwsgi]
command=/home/archi/PycharmProjects/sms/env/bin/uwsgi --http :8000 --module sms.wsgi
directory=/home/archi/PycharmProjects/sms/src

autostart=true
autorestart=true
